# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=su-exec elasticsearch /usr/share/elasticsearch/bin/elasticsearch

[program:kibana]
command=/usr/share/kibana/bin/kibana